# profile/airootfs/etc/systemd/system/getty@tty1.service.d/autologin.conf
# [Service]
# ExecStart=
# ExecStart=-/sbin/agetty -o '-p -f -- \\u' --noclear --autologin root - $TERM

[Service]
ExecStart=
ExecStart=-/sbin/agetty --autologin turtle --noclear %I 38400 linux
